# deploy/pi-gen/stage-ftw/02-wifi-onboarding/files/ftw-wifi-onboarding.service
[Unit]
Description=ftw WiFi onboarding (captive portal fallback)
Documentation=https://github.com/balena-os/wifi-connect
After=NetworkManager.service
Wants=NetworkManager.service
# Skip on subsequent boots once the user has configured WiFi.
ConditionPathExists=!/var/lib/ftw/wifi-configured
# Must run before firstboot — no point docker-pulling over a captive
# portal that hasn't been satisfied yet. firstboot.service is gated
# on network-online.target, which won't flip until NM reports the
# wlan interface as connected, which only happens once the user has
# finished the portal flow.
Before=ftw-firstboot.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/ftw-wifi-onboarding
# wifi-connect drops the AP and exits after the user submits
# credentials. If the whole flow errors out, give the user a fresh
# AP rather than leaving them stuck — but with a cooldown so we
# don't spam-cycle hostapd if something's deeply broken.
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
